# .config/systemd/user/auto-backup.service
[Unit]
Description=Pictures and Videos, ~/.ssh/config to samba

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup.sh
