# etc/systemd/system/smart-test-long.service
[Unit]
Description=hdparm spin down

[Service]
Type=simple
ExecStart=/usr/local/bin/smart_test_long.sh
